# roles/k8s-common/templates/rpfilter-bypass.service.j2
#
# -*- coding:utf-8 -*-
#
# rpfilter-bypass.service
#
# 配置先: /etc/systemd/system/rpfilter-bypass.service
#
# This file is generated by ansible.
# last update: {{ '%Y-%m-%d %H:%M:%S %Z' | strftime(ansible_date_time.epoch) }}
#
[Unit]
Description=Preempt rpfilter drops by accepting return traffic before firewalld
Wants=firewalld.service
After=firewalld.service

[Service]
Type=oneshot
RemainAfterExit=yes
# 起動時, 再読み込み時ともに毎回適用 (firewalld の規則再生成に先行)
ExecStart=/usr/local/sbin/rpfilter-bypass-ensure.sh
ExecReload=/usr/local/sbin/rpfilter-bypass-ensure.sh

[Install]
WantedBy=multi-user.target