# roles/k8s-worker/templates/pin-irqs.service.j2
#
#  -*- coding:utf-8 mode:bash -*-
# This file is generated by ansible.
{# 日付の取得 #}
# last update: {{ '%Y-%m-%d %H:%M:%S %Z' | strftime(ansible_date_time.epoch) }}

#
# /etc/systemd/system/pin-irqs.service
#
[Unit]
Description=Pin all IRQs to housekeeping CPUs
After=multi-user.target

[Service]
Type=oneshot
ExecStart=sudo -- {{k8s_node_setup_tools_dir}}/pin-irqs.py --cpulist {{k8s_reserved_system_cpus_default}} --set-default --set-existing --verbose

[Install]
WantedBy=multi-user.target
