# roles/k8s-worker/templates/pin-worker-queue.service.j2
#
#  -*- coding:utf-8 mode:bash -*-
# This file is generated by ansible.
{# 日付の取得 #}
# last update: {{ '%Y-%m-%d %H:%M:%S %Z' | strftime(ansible_date_time.epoch) }}

#
# /etc/systemd/system/pin-worker-queue.service
#
[Unit]
Description=Pin all unbounded worker queues to system CPUs
After=multi-user.target

[Service]
Type=oneshot
ExecStart={{k8s_node_setup_tools_dir}}/pin-worker-queue.sh

[Install]
WantedBy=multi-user.target
