# roles/user-settings/templates/_config_systemd_user_ssh-agent.service.j2
# -*- mode: systemd -*-
# This file is generated by ansible.
{# 日付の取得 #}
# last update: {{ '%Y-%m-%d %H:%M:%S %Z' | strftime(ansible_date_time.epoch) }}
#
# ssh-agent.service - SSH Agent Service
#
# systemd ユーザーユニット用の ssh-agent サービスファイル。
# ssh-agent.socket からソケットアクティベーションで起動される。
#
# 動作確認コマンド:
#   systemctl --user status ssh-agent.service
#   ssh-add -l

[Unit]
Description=SSH Agent Service
Documentation=man:ssh-agent(1)

[Service]
Type=simple
Environment=SSH_AUTH_SOCK={{ user_settings_systemd_ssh_agent_socket_path }}
ExecStart={{ user_settings_systemd_ssh_agent_bin }} -D -a $SSH_AUTH_SOCK
Restart=on-failure

[Install]
WantedBy=default.target
