# roles/user-settings/templates/_config_systemd_user_ssh-agent.socket.j2
# -*- mode: systemd -*-
# This file is generated by ansible.
{# 日付の取得 #}
# last update: {{ '%Y-%m-%d %H:%M:%S %Z' | strftime(ansible_date_time.epoch) }}
#
# ssh-agent.socket - SSH Agent Socket
#
# systemd ユーザーユニット用の ssh-agent ソケットファイル。
# このソケットを有効化すると, ソケット経由で ssh-agent サービスが自動起動される。
#
# 有効化コマンド:
#   systemctl --user enable ssh-agent.socket
#   systemctl --user start ssh-agent.socket
#
# または:
#   systemctl --user enable --now ssh-agent.socket

[Unit]
Description=SSH Agent Socket

[Socket]
ListenStream={{ user_settings_systemd_ssh_agent_socket_path }}
SocketMode=0600
DirectoryMode=0700

[Install]
WantedBy=sockets.target
